# custom_files/run_bootfile_prod.service
[Unit]
Description = Run standard SLAC boot file
After=network-online.target

[Service]
Type=simple
ExecStart = /root/scripts/run_bootfile_prod.sh

[Install]
WantedBy = default.target
